# config/com.pms.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pms.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/rod/memory/.venv/bin/python</string>
        <string>/Users/rod/memory/scripts/run_sync.py</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>WorkingDirectory</key>
    <string>/Users/rod/memory</string>

    <key>StandardOutPath</key>
    <string>/Users/rod/memory/data/sync_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rod/memory/data/sync_stderr.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
